# old.cron
*/5 * * * * php /home/dennis/noiphp/run.php
